# deploy/configs/supervisor/rush.conf
[program:rush]
command=/var/www/emporium/rush
autostart=true
autorestart=true
user=emporium
