# deploy/configs/supervisor/uwsgi.conf
[program:uwsgi]
command=/var/www/emporium/uwsgi.bin
autostart=true
autorestart=true
user=emporium
stopsignal=QUIT
environment=ENV="PROD"

